# break-reminder.timer
[Unit]
Description=Break reminder every 30 minutes

[Timer]
OnBootSec=30min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
